# macos/LaunchAgents/local.dotfiles.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>local.dotfiles</string>
		<key>ProgramArguments</key>
		<array>
			<string>/Users/ericdahlseng/.dotfiles/bin/dot</string>
			<string>-g</string>
			<string>-d</string>
			<string>-i</string>
			<string>-l</string>
		</array>
		<key>StartCalendarInterval</key>
		<dict>
			<key>Hour</key>
			<integer>8</integer>
			<key>Minute</key>
			<integer>30</integer>
		</dict>
	</dict>
</plist>
